# packer/files/etc/init/zuul-init.conf
description "zuul"

start on started elastic-network-interfaces
stop on run level [!2345]
env enabled=1


exec  /usr/local/bin/zuul-init.sh
